# etc/logrotate.d/rc.local
/var/log/rc.local.log {
	daily
	missingok
	rotate 14
	compress
	delaycompress
	notifempty
}
