# profile.d/git-credential-manager.sh

if [ -x /usr/local/bin/git-credential-manager ]
then
    export GCM_CREDENTIAL_STORE="cache"
fi
