# profile.d/teleport.sh

export TELEPORT_PROXY=teleport.jmsh-home.com:444
alias tsh-ssh='tsh ssh -A'
