# cron.daily/usr-local-scripts.sh
#!/bin/bash

# cron.daily script to update local scripts

exec /usr/local/scripts/install.sh
